# app/server/hamster-app.service
[Unit]
Description=hamster-cam Fastify + tRPC app server
Documentation=https://github.com/aaron/hamster-cam
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Default service account. Do NOT hardcode your own user here — override
# per host with a drop-in instead (never committed, survives deploys):
#   /etc/systemd/system/hamster-app.service.d/override.conf
#   [Service]
#   User=<your-account>
#   Group=<your-account>
# See docs/SETUP_MAC_MINI.md step 10.3.
User=hamster
Group=hamster
WorkingDirectory=/opt/hamster-cam/app/server
EnvironmentFile=/opt/hamster-cam/.env
ExecStart=/usr/bin/node /opt/hamster-cam/app/server/dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Sensible hardening — tweak per deployment if it interferes with /opt/hamster-cam.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
