# pi-zero/go2rtc.service
[Unit]
Description=go2rtc — RTSP listener for the hamster-cam Pi Zero
Documentation=https://github.com/AlexxIT/go2rtc https://github.com/aaroncoppock/hamster-cam
# ntp-sync.service gates on a confirmed NTP synchronization so the first
# frame timestamp is correct. Without this, a Pi that lost power for hours
# starts go2rtc with a clock that is hours behind the Mac Mini, producing
# negative zone-visit durations that the diary pipeline silently drops.
After=network-online.target ntp-sync.service
Wants=network-online.target ntp-sync.service

[Service]
Type=simple
# RTSP_PASSWORD ships in this root-only env file (chmod 600). Keeping the
# secret here rather than on the ExecStart command line keeps it out of
# /proc/<pid>/cmdline and out of `ps` output.
EnvironmentFile=/etc/go2rtc/go2rtc.env
# Anti-flicker / anti-banding: US mains run at 60 Hz, so the sensor's
# electronic rolling shutter beats against the light cycle and produces
# horizontal wavy "heat-haze" distortion under artificial light.
# power_line_frequency=2 → 60 Hz; =1 → 50 Hz (UK/EU). This must be set
# BEFORE go2rtc/ffmpeg opens /dev/video0, which is why it lives in
# ExecStartPre rather than a one-shot post-start. The leading `-` means
# "ignore a non-zero exit" — if this Pi's sensor driver doesn't expose the
# control (e.g. a substitute camera module) the failure is absorbed and
# go2rtc still starts. Path /usr/bin/v4l2-ctl is from the Debian bookworm
# v4l-utils package; run `which v4l2-ctl` on the Pi if a future OS upgrade
# moves it. Re-applied on every start, so it survives the watchdog restarts
# and the documented cam2 random reboots.
ExecStartPre=-/usr/bin/v4l2-ctl -d /dev/video0 --set-ctrl=power_line_frequency=2
ExecStart=/usr/local/bin/go2rtc -config /etc/go2rtc/go2rtc.yaml
Restart=always
RestartSec=5

# User=root is intentional. The IMX462 lives at /dev/video0, which is owned
# root:video on Raspberry Pi OS. The cleaner alternative is to run as a
# non-root user added to the `video` group — but that user also needs to
# bind :1984 and :8554 (both >1024, so unprivileged is fine there). If you
# tighten this, also drop the EnvironmentFile permissions to match the new
# uid. Keeping root for the family-cam deployment because the watchdog
# (which calls `systemctl restart go2rtc` and `systemctl reboot`) is itself
# a root unit, and one privilege boundary is simpler to reason about.
User=root
Group=root

# Defensive hardening — go2rtc only needs to read /dev/video0, write to its
# tmp dir, and talk to the network.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log /run

[Install]
WantedBy=multi-user.target
